# services/smart_speaker.service
[Unit]
Description=Smart Speaker IoT Service
After=network-online.target mopidy.service
Wants=network-online.target

[Service]
Type=simple
User=iot-proj
WorkingDirectory=/home/iot-proj/IOT-project--Smart-Speaker/Main
ExecStart=/home/iot-proj/IOT-project--Smart-Speaker/Main/main.py
Restart=on-failure
RestartSec=10

# Logging - output to file for debug endpoints + journal for systemctl
# View with: journalctl -u smart_speaker -f  OR  tail -f /var/log/smart_speaker.log
StandardOutput=append:/var/log/smart_speaker.log
StandardError=append:/var/log/smart_speaker.log

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
